# charmtools/templates/python_services/files/templates/upstart.conf
description "$metadata.package"
author "$metadata.maintainer"

start on runlevel [2345]
stop on runlevel [016]

respawn

console log
script
    echo Fake service; sleeping for an hour...
    sleep 360
end script
